# infinitive/rootfs/etc/services.d/infinitive/run
#!/usr/bin/env bashio
# ==============================================================================
# Start the infinitive service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare serial_tty
declare serial_addr
declare serial_port
declare serialPort

## Get the 'message' key from the user config options.
serial_tty=$(bashio::config 'serial_tty')
serial_addr=$(bashio::config 'serial_addr')
serial_port=$(bashio::config 'serial_port')

if [ ! -z "$serial_addr" ]; then
  serialPort="tcp://$serial_addr:$serial_port"
elif [ ! -z "$serial_tty" ]; then
  serialPort="$serial_tty"
else
  bashio::log.error "Serial port not defined!"
  exit 1
fi

## Print the serial port we're using
bashio::log.info "Using serial ${serialPort}"

## Run your program
exec /usr/bin/infinitive -serial "$serialPort"
