# infinitive/rootfs/etc/services.d/infinitive/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the infinitive service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare serial

## Get the 'message' key from the user config options.
serial=$(bashio::config 'serial')

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "${serial:="/dev/ttyS0"}"

## Run your program
exec /usr/bin/infinitive -serial $serial
